# conf/saltstack_minion_supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service sshd start
autorestart=false

[program:saltstack]
command=service salt-minion start
autorestart=false
