# conf/saltstack_master_supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service sshd start
autorestart=false

[program:saltstack]
command=service salt-master start
autorestart=false
